# network_monitor/rootfs/etc/s6-overlay/s6-rc.d/network-monitor/run
#!/command/with-contenv bashio

CONFIG=/data/options.json

export NETWORK_RANGE=$(bashio::config 'network_range')
export SCAN_INTERVAL_HOURS=$(bashio::config 'scan_interval_hours')
export OFFLINE_THRESHOLD_HOURS=$(bashio::config 'offline_threshold_hours')
export TELEGRAM_ENABLED=$(bashio::config 'telegram_enabled')
export TELEGRAM_BOT_TOKEN=$(bashio::config 'telegram_bot_token')
export TELEGRAM_CHAT_ID=$(bashio::config 'telegram_chat_id')
export API_KEY=$(bashio::config 'api_key')
export FLASK_PORT=5004
export INGRESS_PATH=$(bashio::addon.ingress_entry)

bashio::log.info "Starting Network Monitor..."
bashio::log.info "Network range: ${NETWORK_RANGE}"
bashio::log.info "Scan interval: ${SCAN_INTERVAL_HOURS}h"
bashio::log.info "Ingress path: ${INGRESS_PATH}"

cd /app && exec python3 wsgi.py
