# network_monitor/rootfs/etc/services.d/network-monitor/run
#!/bin/bash
set -e

CONFIG=/data/options.json

export NETWORK_RANGE=$(jq --raw-output '.network_range // "192.168.1.0/24"' "${CONFIG}")
export SCAN_INTERVAL_HOURS=$(jq --raw-output '.scan_interval_hours // 1' "${CONFIG}")
export OFFLINE_THRESHOLD_HOURS=$(jq --raw-output '.offline_threshold_hours // 2' "${CONFIG}")
export TELEGRAM_BOT_TOKEN=$(jq --raw-output '.telegram_bot_token // ""' "${CONFIG}")
export TELEGRAM_CHAT_ID=$(jq --raw-output '.telegram_chat_id // ""' "${CONFIG}")
export API_KEY=$(jq --raw-output '.api_key // ""' "${CONFIG}")
export FLASK_PORT=5000
export DEBUG=false

echo "[network-monitor] Starting..."
echo "[network-monitor] Network: ${NETWORK_RANGE}"
echo "[network-monitor] Scan every: ${SCAN_INTERVAL_HOURS}h"

cd /app && exec python3 main.py
